# includes/supervisor/conf.d/gunicorn.conf
[program:gunicorn]
command = /usr/local/bin/gunicorn -b '0.0.0.0:8000'
stdout_events_enabled = true
stderr_events_enabled = true
autorestart = true
user = zulip
group = zulip
